# deploy/supervisord.conf
[program:hello]
command = /home/phodal/growth-studio/deploy/gunicorn_start
stdout_logfile = /home/phodal/growth-studio/logs/gunicorn_supervisor.log
redirect_stderr = true
